# src/scripts/nullsd-init.sh.in
#!/bin/sh
#
# Init script to start/stop/restart the NullLogic GroupServer servers.
# On SYSV init systems, this script is usually called 'nullsd',
# and located in '/etc/init.d'.
#
# chkconfig: 2345 80 10
# description: NullLogic GroupServer Server.
# processname: nullsd
# pidfile: /var/run/nullsd.pid

# Source function library.
if [ -x /etc/rc.d/init.d/functions ]; then
	. /etc/rc.d/init.d/functions
fi
# Source networking configuration.
if [ -x /etc/sysconfig/network ]; then
	. /etc/sysconfig/network
fi
# Include nullsd defaults if available
if [ -f /etc/default/nullsd ]; then
	. /etc/default/nullsd
fi

prefix=@prefix@
exec_prefix=@exec_prefix@
BINDIR="@bindir@"
ETCDIR="@sysconfdir@"
LIBDIR="@libdir@"
VARDIR="@localstatedir@"
PATH=$BINDIR:$PATH
LD_LIBRARY_PATH=$LIBDIR:$LD_LIBRARY_PATH
pidfile="/var/run/nullsd.pid"

check_status() {
	if [ -f ${pidfile} ]; then
		if [ -f "/proc/`cat ${pidfile}`/status" ]; then
			SERVER_STATUS="running"
		else
			SERVER_STATUS="not running"
		fi
	else
		SERVER_STATUS="not running"
	fi
#	RC=`lsof -t -a -i :8086 -i :8112 -sTCP:LISTEN`
#	RC=`ps ax | grep -v grep | grep -c [n]ullsd`
#	if test ${RC} -gt 0; then
}

nullsd_start() {
	if [ "$SERVER_STATUS" = "running" ]; then
		echo "NullLogic GroupServer appears to be running already."
	else
		echo -n "Starting NullLogic GroupServer services: "
		cd ${BINDIR}
# You can either load all the servers as one process,
		echo -n "OK" ; nullsd
# or you can load them separately.
#		echo -n "httpd " ; nullsd httpd
#		echo -n "pop3d " ; nullsd pop3d
#		echo -n "smtpd " ; nullsd smtpd
#		echo -n "smtpq " ; nullsd smtpq
		echo ""
	fi
}

nullsd_stop() {
	if [ "$SERVER_STATUS" = "running" ]; then
		echo "Stopping NullLogic GroupServer services"
		if [ -f /var/run/nullsd.pid ]; then
			kill -15 `cat /var/run/nullsd.pid`
			sleep 2
			kill -9 `cat /var/run/nullsd.pid` 1>/dev/null 2>&1
			rm -f /var/run/nullsd.pid
		else
			killall -15 nullsd 2>/dev/null
			sleep 2
			killall -9 nullsd 2>/dev/null
		fi
	else
		echo "NullLogic GroupServer does not appear to be running."
	fi
}

nullsd_status() {
	if [ "$SERVER_STATUS" = "running" ]; then
		echo "NullLogic GroupServer appears to be running"
	else
		echo "NullLogic GroupServer does not appear to be running."
	fi
}

nullsd_restart() {
	nullsd_stop
	check_status
	nullsd_start
}

check_status
case "$1" in
start)
	nullsd_start
	;;
stop)
	nullsd_stop
	exit 0
	;;
status)
	nullsd_status
	;;
restart|reload)
	nullsd_restart
	;;
*)
	echo "Usage: $0 {start|stop|status|restart|reload}"
esac

exit $?
