# src/scripts/nullgs-init.sh.in
#!/bin/sh
#
# Init script to start/stop/restart the NullLogic GroupServer servers.
# On SYSV init systems, this script is usually called 'nullgs',
# and located in '/etc/init.d'.
#
# chkconfig: 2345 80 10
# description: NullLogic GroupServer Server.
# processname: nullgsd
# pidfile: /var/run/nullgsd.pid

# Source function library.
if [ -x /etc/rc.d/init.d/functions ]; then
	. /etc/rc.d/init.d/functions
fi
# Source networking configuration.
if [ -x /etc/sysconfig/network ]; then
	. /etc/sysconfig/network
fi
# Include nullgs defaults if available
if [ -f /etc/default/nullgs ]; then
	. /etc/default/nullgs
fi

prefix=@prefix@
exec_prefix=@exec_prefix@
BINDIR="@bindir@"
ETCDIR="@sysconfdir@"
LIBDIR="@libdir@"
VARDIR="@localstatedir@"
PATH=$BINDIR:$PATH
LD_LIBRARY_PATH=$LIBDIR:$LD_LIBRARY_PATH

check_status() {
	RC=`ps ax | grep -v grep | grep -c [n]ullgsd`
	if test ${RC} -gt 0; then
		SERVER_STATUS="running"
	else
		SERVER_STATUS="not running"
	fi
}

nullgs_start() {
	if [ "$SERVER_STATUS" = "running" ]; then
		echo "NullLogic GroupServer appears to be running already."
	else
		echo -n "Starting NullLogic GroupServer services: "
		cd ${BINDIR}
# You can either load all the servers as one process,
		echo -n "OK" ; nullgsd
# or you can load them separately.
#		echo -n "httpd " ; nullgsd httpd
#		echo -n "pop3d " ; nullgsd pop3d
#		echo -n "smtpd " ; nullgsd smtpd
#		echo -n "smtpq " ; nullgsd smtpq
		echo ""
	fi
}

nullgs_stop() {
	if [ "$SERVER_STATUS" = "running" ]; then
		echo "Stopping NullLogic GroupServer services"
		if [ -f /var/run/nullgsd.pid ]; then
			kill -15 `cat /var/run/nullgsd.pid`
			sleep 2
			kill -9 `cat /var/run/nullgsd.pid` 1>/dev/null 2>&1
			rm -f /var/run/nullgsd.pid
		else
			killall -15 nullgsd 2>/dev/null
			sleep 2
			killall -9 nullgsd 2>/dev/null
		fi
	else
		echo "NullLogic GroupServer does not appear to be running."
	fi
}

nullgs_status() {
	if [ "$SERVER_STATUS" = "running" ]; then
		echo "NullLogic GroupServer appears to be running"
	else
		echo "NullLogic GroupServer does not appear to be running."
	fi
}

nullgs_restart() {
	nullgs_stop
	check_status
	nullgs_start
}

check_status
case "$1" in
start)
	nullgs_start
	;;
stop)
	nullgs_stop
	exit 0
	;;
status)
	nullgs_status
	;;
restart|reload)
	nullgs_restart
	;;
*)
	echo "Usage: $0 {start|stop|status|restart|reload}"
esac

exit $?
